# src/DigitalSignage.Client.RaspberryPi/digitalsignage-client.service
[Unit]
Description=Digital Signage Client
After=network-online.target graphical.target multi-user.target
Wants=network-online.target

[Service]
Type=simple
User=INSTALL_USER
Group=INSTALL_USER
WorkingDirectory=/opt/digitalsignage-client
Environment="PYTHONUNBUFFERED=1"
Environment="QT_QPA_PLATFORM=xcb"
Environment="QT_LOGGING_RULES=*.debug=false"
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/INSTALL_USER/.Xauthority"

# Pre-start file checks
ExecStartPre=/bin/bash -c 'test -f /opt/digitalsignage-client/start-with-display.sh || (echo "ERROR: start-with-display.sh not found" && exit 1)'
ExecStartPre=/bin/bash -c 'test -x /opt/digitalsignage-client/start-with-display.sh || (echo "ERROR: start-with-display.sh not executable" && exit 1)'
ExecStartPre=/bin/bash -c 'test -f /opt/digitalsignage-client/venv/bin/python3 || (echo "ERROR: Python venv not found" && exit 1)'

# CRITICAL FIX: Start X11 without desktop environment if not running (for HDMI displays)
# This prevents desktop terminal window from appearing over Digital Signage client
# Two scenarios:
#   1. X11 already running on :0 (desktop booted) → Use it directly with start-with-display.sh
#   2. X11 not running (console boot) → Start X11 via startx with .xinitrc (no desktop)
# The .xinitrc will execute start-with-display.sh as the ONLY GUI application
ExecStart=/bin/bash -c 'if DISPLAY=:0 xset q &>/dev/null 2>&1; then echo "X11 detected on :0, using existing display"; exec /opt/digitalsignage-client/start-with-display.sh; else echo "No X11 detected, starting X11 without desktop via startx"; exec /usr/bin/startx /home/INSTALL_USER/.xinitrc -- :0 vt1 -nocursor; fi'

# Cleanup
ExecStopPost=/bin/bash -c 'pkill -f "Xvfb :99" || true'

# Restart configuration
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=digitalsignage-client

# Security
# NoNewPrivileges=true  # DISABLED: Prevents systemctl reboot command from working
PrivateTmp=true

# Restart policy
StartLimitInterval=200
StartLimitBurst=5

[Install]
WantedBy=graphical.target
