# src/DigitalSignage.Client.RaspberryPi/digitalsignage-client.service
[Unit]
Description=Digital Signage Client
After=network-online.target graphical.target
Wants=network-online.target

[Service]
Type=notify
NotifyAccess=main
User=INSTALL_USER
Group=INSTALL_USER
WorkingDirectory=/opt/digitalsignage-client
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/INSTALL_USER/.Xauthority"
ExecStartPre=/bin/sleep 10
ExecStart=/opt/digitalsignage-client/venv/bin/python3 /opt/digitalsignage-client/client.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=digitalsignage-client

# Watchdog - will restart if no ping received within 60 seconds
WatchdogSec=60
# Give service 30 seconds to start and notify ready
TimeoutStartSec=30

# Security
NoNewPrivileges=true
PrivateTmp=true

# Restart policy
StartLimitInterval=200
StartLimitBurst=5

[Install]
WantedBy=graphical.target
